# services/digital-signage.service
[Unit]
Description=SignageCommander Digital Signage Platform
Documentation=https://github.com/signagecommander/platform
After=network.target network-online.target
Wants=network-online.target
RequiresMountsFor=/opt/digital-signage

[Service]
Type=simple
User=signage
Group=signage
WorkingDirectory=/opt/digital-signage
ExecStart=/usr/bin/python3 /opt/digital-signage/main.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=digital-signage

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/opt/digital-signage
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=PYTHONPATH=/opt/digital-signage
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target